# pacrunner/src/dbus-pacrunner.service.in
[Unit]
Name=Pacrunner service

[Service]
Type=dbus
BusName=org.pacrunner
ExecStart=@prefix@/sbin/pacrunner

